# ansible/tasks/app/docker/containers/promgraf/templates/node_exporter.service.j2
[Unit]
Description=starts node-exporter service

[Service]
ExecStart=/usr/bin/node_exporter
Type=simple
Restart=on-failure

[Install]
WantedBy=multi-user.target